# net-wireless/aic8800-driver/files/init/load-aic_wifi-modules.conf
#  Copyright 2024 The FydeOS Authors. All rights reserved.
#  Use of this source code is governed by a BSD-style license that can be
#  found in the LICENSE file.
#  Author: Yang Tsao<yang@fydeos.io>

description "Job to load modules"
author      "yang@fydeos.io"

oom score never
task

start on login-prompt-visible

script
  modprobe -q -a sdhci-rtk rtk_rfkill
  rfkill unblock 0
  sleep 1
  modprobe -q aic8800_fdrv
  hciattach /dev/ttyS1 any 1500000 flow
end script
